# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=Xvfb :1 -screen 0 1080x1920x24
autorestart=true

[program:fluxbox]
command=fluxbox
environment=DISPLAY=:1
autorestart=true

[program:emulator]
command=emulator -avd MyCloudPhone -gpu swiftshader_indirect -no-audio -qemu -m 2048
environment=DISPLAY=:1
autorestart=true

[program:x11vnc]
command=x11vnc -display :1 -nopw -forever -shared -rfbport 5900
autorestart=true

[program:novnc]
command=/usr/share/novnc/utils/launch.sh --vnc localhost:5900 --listen 6080
autorestart=true
